# conf/systemd.service
[Unit]
Description=NodeBB
Documentation=https://docs.nodebb.org
After=system.slice multi-user.target mongod.service

[Service]
Type=forking
User=__APP__
Group=__APP__
Environment=NODE_ENV=production
WorkingDirectory=__FINALPATH__
PIDFile=__FINALPATH__/pidfile
Environment="__YNH_NODE_LOAD_PATH__"
ExecStart=__YNH_NODE__ loader.js
Restart=always
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=__APP__

[Install]
WantedBy=multi-user.target
